# novnc_rocker/templates/turbovnc.conf
[program:turbovnc]
command=/opt/TurboVNC/bin/vncserver  -SecurityTypes None -fg -nohttpd
autorestart=true