# novnc_rocker/templates/turbovnc.conf
[program:turbovnc]
command=/opt/TurboVNC/bin/vncserver  -SecurityTypes None -fg -nohttpd -geometry 1920x1080
autorestart=true